# docker/lib/base_tomcat_image/container-files/etc/systemd/system/tomcat_image_init.service
[Unit]
Description=Used with systemd images to run stuff at startup
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tomcat_image_init.sh
TimeoutStartSec=25

[Install]
WantedBy=default.target