# deploy/lib/vagrant/etc/profile.d/custom.sh
export WGET_CACHE_DIR=/opt/wget_cache
export PGDATA=/shared/pg_data

export JAVA_HOME=/usr/lib/jvm/jre

export RD_PASSWORD=admin
export RD_URL=http://127.0.0.1:4440
export RD_USER=admin
export RD_PASSWORD=admin